# pritunl-link/pritunl-link.service
[Unit]
Description=Pritunl Link Daemon

[Service]
LimitNOFILE=50000
ExecStart=/usr/bin/pritunl-link start

[Install]
WantedBy=multi-user.target
